# meli.service
[Unit]
Description=Meli Honeypot Command Center
After=graphical-session.target meli-ingest.service
Wants=meli-ingest.service

[Service]
Type=simple
ExecStart=/usr/local/bin/meli
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0
Environment=XDG_RUNTIME_DIR=/run/user/%U

[Install]
WantedBy=graphical-session.target
