# meli-labyrinth-digest.service
[Unit]
Description=Meli Labyrinth daily digest builder
Documentation=https://github.com/sierengowskisierengowski-cpu/Meli
After=meli.service
Wants=meli.service

[Service]
Type=oneshot
ExecStart=/opt/meli/venv/bin/python -m meli.labyrinth.digest --hours 24
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
# Hardening — runs read-only against replay logs, write-only to its
# own digest output directory.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/share/meli

[Install]
WantedBy=default.target
