# meli-labyrinth-digest.timer
[Unit]
Description=Meli Labyrinth daily digest — fires once a day at 07:00 local
Documentation=https://github.com/sierengowskisierengowski-cpu/Meli

[Timer]
# Fires at 07:00 every day. RandomizedDelaySec spreads load if multiple
# Meli instances run on the same host (rare but harmless).
OnCalendar=*-*-* 07:00:00
RandomizedDelaySec=5min
# Run on next boot if the system was off at 07:00 — operator should
# always get yesterday's digest even if the box was rebooting.
Persistent=true

[Install]
WantedBy=timers.target
