# meli-ingest.service
[Unit]
Description=Meli Ingest Daemon — Background honeypot event ingestion
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
ExecStart=/usr/local/bin/meli --daemon ingest
Restart=always
RestartSec=10
StandardOutput=append:%h/.local/share/meli/logs/ingest.log
StandardError=append:%h/.local/share/meli/logs/ingest.log

[Install]
WantedBy=default.target
